# analyze/stampedeLariatReport.cron
#!/bin/bash
# -*- shell-script -*-

today=$(date +%F)
fn30="LariatReport-last30-$today.txt"


/scratch/projects/lariatData/analyze_LDDjob $(find . -name 'lariatData-sgeT*.lua' -mmin -43200 ) > last30/$fn30

fn365="LariatReport-last365-$today.txt"

/scratch/projects/lariatData/analyze_LDDjob $(find . -name 'lariatData-sgeT*.lua' -mmin -525600 ) > last365/$fn365
